# SockBot.conf
start on local-filesystem
stop on shutdown

console log

chdir /home/accalia/SockBotProd

script
    exec su accalia -c "node SockBot.js users/SockBot.json"
end script
